# vps/etc/systemd/system/dotsmith.service
[Unit]
Description=dotsmith orchestrator (FastAPI + Vite)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=dotsmith
WorkingDirectory=/srv/dotsmith/backend
ExecStart=/srv/dotsmith/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
