# vps/etc/systemd/system/dotsmith-pixellab-snapshot.service
[Unit]
Description=dotsmith — snapshot do saldo das contas Pixellab
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=dotsmith
Group=dotsmith
WorkingDirectory=/srv/dotsmith/backend
ExecStart=/srv/dotsmith/backend/.venv/bin/python -m app.jobs.snapshot_balances
