# vps/etc/systemd/system/lgmateus.service
[Unit]
Description=lgmateus.com (Next.js)
After=network.target

[Service]
Type=simple
User=lgmateus
WorkingDirectory=/srv/lgmateus
Environment=NODE_ENV=production
Environment=PORT=3000
# mise exec resolve o node do usuario lgmateus (sem path hardcoded de versao)
ExecStart=/srv/lgmateus/.local/bin/mise exec -- node node_modules/next/dist/bin/next start --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
